# distributable/tiedot.service
[Unit]
Description=tiedot database server (HTTP API service)
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/tiedot
ExecStart=/usr/bin/tiedot -mode=httpd -port=${HTTP_PORT} -dir=${HTTP_DB_DIR} ${HTTP_EXTRA_ARGS}
ExecStop=/bin/sh -c 'curl -s "http://localhost:${HTTP_PORT}/shutdown" > /dev/null 2>&1 || true'

[Install]
WantedBy=multi-user.target

